# templates/assistant-ssh-gateway.service
[Unit]
Description=Rebuild SSH authorized_keys from assistant profiles

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rebuild-authorized-keys
