# templates/domovoy-ssh-gateway.service
[Unit]
Description=Rebuild SSH authorized_keys from domovoy profiles

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rebuild-authorized-keys
