# templates/domovoy-ssh-gateway.path
[Unit]
Description=Watch domovoy profiles for new SSH keys

[Path]
PathChanged=%h/setup

[Install]
WantedBy=default.target
